# oalpr_local/rootfs/etc/services.d/alpr/run

#!/usr/bin/with-contenv bash
# s6 запускает этот скрипт как главный процесс аддона
exec /usr/bin/python3 -u /opt/server.py
