# oalpr_local/rootfs/etc/services.d/alpr/run
#!/usr/bin/with-contenv bash
exec /usr/bin/python3 -u /opt/server.py
